# deploy/supervisor/hotel-pms-worker.conf
; ==================================================
; Supervisor Config — Hotel PMS OTA Autopilot
; ==================================================
; Copy this file to: /etc/supervisor/conf.d/hotel-pms-worker.conf
;
; Then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start hotel-pms-worker:*
;   sudo supervisorctl start hotel-pms-scheduler:*
;
; Check status:
;   sudo supervisorctl status
;
; View logs:
;   sudo tail -f /var/log/hotel-pms-worker.log
;   sudo tail -f /var/log/hotel-pms-scheduler.log
; ==================================================

[program:hotel-pms-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hotel-pms/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --backoff=30
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/hotel-pms-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www-data
directory=/var/www/hotel-pms
environment=APP_ENV="production",APP_DEBUG="false"
stopsignal=SIGTERM
stopwaitsecs=30

[program:hotel-pms-scheduler]
process_name=%(program_name)s
command=php /var/www/hotel-pms/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/hotel-pms-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
user=www-data
directory=/var/www/hotel-pms
environment=APP_ENV="production",APP_DEBUG="false"
stopsignal=SIGTERM
stopwaitsecs=10
